# launchd/com.tubi.babysit-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Optional: run the babysit-agent daemon at login.
  1. Edit the paths below (NODE_BIN, WorkingDirectory, log paths).
  2. cp this file to ~/Library/LaunchAgents/
  3. launchctl load ~/Library/LaunchAgents/com.tubi.babysit-agent.plist
  Stop with: launchctl unload ~/Library/LaunchAgents/com.tubi.babysit-agent.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tubi.babysit-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/leo/.nvm/versions/node/v22.13.0/bin/node</string>
    <string>/Users/leo/dev/tubi/ai_workspaces/babysit-agent/packages/server/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/leo/dev/tubi/ai_workspaces/babysit-agent</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/leo/.babysit-agent/daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/leo/.babysit-agent/daemon.err.log</string>
</dict>
</plist>
