# launchd/io.tubi.babysit-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.tubi.babysit-agent</string>

    <!-- Run `npm run dev:server` (tsx watch) via the nvm node that provides npm. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/leo/.nvm/versions/node/v22.13.0/bin/node</string>
        <string>/Users/leo/.nvm/versions/node/v22.13.0/lib/node_modules/npm/bin/npm-cli.js</string>
        <string>run</string>
        <string>dev:server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/leo/dev/tubi/ai_workspaces/babysit-agent</string>

    <!-- launchd starts with a bare PATH; the daemon shells out to `gh` (Homebrew)
         and spawns `node` for tsx, so both dirs must be present. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/leo/.nvm/versions/node/v22.13.0/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/leo</string>
        <!-- Root all heavy state (SQLite + repo clones + worktrees + CI logs)
             under the workspace `.data/` — the SAME dir the Docker path
             bind-mounts. A native process reads it at full disk speed (none of
             the VirtioFS tax the container pays), so switching between the
             launchd and Docker paths reuses one state dir: the www clone and
             any in-flight threads survive the switch. -->
        <key>BABYSIT_DATA_DIR</key>
        <string>/Users/leo/dev/tubi/ai_workspaces/babysit-agent/.data</string>
    </dict>

    <!-- Login item: start at load and keep it running. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Don't hammer restart-loop if it crashes on startup. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logs live under the same `.data/` as state — one dir to tail. -->
    <key>StandardOutPath</key>
    <string>/Users/leo/dev/tubi/ai_workspaces/babysit-agent/.data/daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leo/dev/tubi/ai_workspaces/babysit-agent/.data/daemon.err.log</string>
</dict>
</plist>
